# pi-mqtt-server/fish-feeder.service
[Unit]
Description=Fish Feeder IoT System Auto-Start
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/fish-feeder/pi-mqtt-server
ExecStart=/usr/bin/python3 /home/pi/fish-feeder/pi-mqtt-server/autostart.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/home/pi/fish-feeder/pi-mqtt-server
Environment=PATH=/usr/local/bin:/usr/bin:/bin

# Resource limits
LimitNOFILE=65536
MemoryLimit=512M

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target 